# terraform/modules/clusters/cluster/instance_group/templates/k8s-vip.service.tmpl
[Unit]
Description=Kubernetes VIP Healthcheck (control-plane only)
After=network-online.target frr.service
Wants=network-online.target
Requires=frr.service
# Only start on nodes that actually run the apiserver static pod
ConditionPathExists=/etc/kubernetes/manifests/kube-apiserver.yaml

[Service]
Type=oneshot
Environment=RISE=${rise}
Environment=FALL=${fall}
Environment=COOLDOWN=${cooldown}
ExecStart=/usr/local/bin/k8s-vip-healthcheck.sh
User=root
Group=root

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/run/k8s-vip
PrivateTmp=true
ProtectKernelTunables=false
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
TimeoutStartSec=30
TimeoutStopSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=k8s-vip